# infrastructure/ansible/roles/notification-service/templates/notification-service.service.j2
[Unit]
Description=Notification Service
After=network.target

[Service]
EnvironmentFile=/etc/notification-service.env
ExecStart=/usr/bin/java -jar /opt/notification-service/notification-service.jar
User=root
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
